# etc/init.d/waydroid-openrc
#!/sbin/openrc-run

description="Manages the Waydroid container"

depend() {
	need net
}

start_pre() {
	modprobe -q loop
	modprobe -q fuse
}

start() {
	ebegin "Starting Waydroid container"
	/usr/bin/waydroid -w container start
	eend $?
}

stop() {
	ebegin "Stopping Waydroid session and container"
	/usr/bin/waydroid session stop
	/usr/bin/waydroid container stop
	eend 0
}
