# debian/postinst
#! /bin/bash

#   Add waydroid service files to runlevel/default

rc-update add waydroid-openrc default

#   Update services list

rc-update -u
